# docker/supervisor/0-redis.conf
[program:redis]
command=redis-server
autostart=true
autorestart=true
directory=/var/lib/redis
priority=1
