# docker/supervisor/1-nginx.conf
[program:nginx]
command = /usr/sbin/nginx
autostart = true
autorestart=true
priority=2
