# docker/supervisor/2-php.conf
[program:php]
command = php-fpm -F
autostart=true
autorestart=true
priority=3
