# docker/supervisor/4-worker.conf
[program:pkg-workers]
environment =
        HOME=/var/www/
command=/var/www/packagist/bin/console packagist:run-workers -v --env=prod --no-debug
directory=/var/www/packagist/
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
startsecs=0
stdout_logfile=/proc/self/fd/2
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=5
user=www-data
